# libs/openldap/files/ldap.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2009-2011 OpenWrt.org

START=61

SERVICE_USE_PID=1
DATADIR=/etc/openldap/data
PIDFILE=/var/run/slapd.pid

start() {
	mkdir -m 0755 -p $DATADIR
	service_start /usr/sbin/slapd
}

stop() {
	service_stop /usr/sbin/slapd
}

reload() {
	service_reload /usr/sbin/slapd
}
